# core/templates/patrol.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  QA patrol beacon. Installed by
  `scripts/install.sh --enable-auto-patrol` at project setup time (opt-in;
  default install does NOT install this plist — patrol runs on operator
  natural-language request instead).

  Substitutions performed at install time:
    {PROJECT}          → project name (e.g. install)
    {CADENCE_SECONDS}  → patrol cadence, seconds (install default 1800 / 30min, overridable)
    {CLAWSEAT_ROOT}    → absolute path to ClawSeat (e.g. <HOME>/.clawseat)
    {LOG_DIR}          → directory for stdout/err logs (e.g. ~/.agents/tasks/<project>/patrol/logs)

  Install location:
    ~/Library/LaunchAgents/com.clawseat.{PROJECT}.patrol.plist

  Activation:
    launchctl bootstrap gui/$(id -u) <plist-path>
    launchctl kickstart  gui/$(id -u)/com.clawseat.{PROJECT}.patrol

  Deactivation (e.g. on `agent-admin project remove`):
    launchctl bootout    gui/$(id -u)/com.clawseat.{PROJECT}.patrol

  Behavior: every CADENCE_SECONDS, this plist runs the one-liner below. The
  one-liner uses the canonical send-and-verify wrapper to send a **natural
  language** patrol request to the memory tmux session if that
  session currently exists. Natural language (not a `/`-prefixed slash
  command) is required because Claude Code treats `/xxx` as a slash-command
  resolver lookup and returns "Unknown command" for un-registered tokens
  like `/patrol-tick`. Natural language is identified by memory SKILL §3
  semantic-trigger detection.

  If the session is gone (user killed it or hasn't started it yet), the
  send is silently skipped so we don't error-spam the logs.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clawseat.{PROJECT}.patrol</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>SESSION="$('{CLAWSEAT_ROOT}/core/shell-scripts/agentctl.sh' session-name memory --project '{PROJECT}' 2&gt;/dev/null || true)"; [[ -n "$SESSION" ]] &amp;&amp; tmux has-session -t "=$SESSION" 2&gt;/dev/null &amp;&amp; '{CLAWSEAT_ROOT}/core/shell-scripts/send-and-verify.sh' --project '{PROJECT}' memory '请按 SKILL §3 做一次 Phase-B 稳态巡检（P1..P7）。Please run one Phase-B patrol cycle per SKILL §3.' &gt;/dev/null 2&gt;&amp;1 || true</string>
    </array>

    <key>StartInterval</key>
    <integer>{CADENCE_SECONDS}</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>{LOG_DIR}/patrol.out.log</string>

    <key>StandardErrorPath</key>
    <string>{LOG_DIR}/patrol.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
