# scripts/com.clawseat.task-watch.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Template — render via: sed "s|{CLAWSEAT_ROOT}|$CLAWSEAT_ROOT|g" com.clawseat.task-watch.plist.in > com.clawseat.task-watch.plist -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.clawseat.task-watch</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>{CLAWSEAT_ROOT}/scripts/task-watch.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>StandardOutPath</key>
  <string>{CLAWSEAT_ROOT}/.agent/task-watch/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>{CLAWSEAT_ROOT}/.agent/task-watch/launchd.err.log</string>
</dict>
</plist>
